# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/timeserver/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; -----------------------------------------------------------------------
; 1. chrony — startuje pierwszy, tworzy segment SHM 0x4E545030 (NTP0)
;    gps-time-bridge.py oraz server.py używają tego segmentu
; -----------------------------------------------------------------------
[program:chrony]
command=/usr/sbin/chronyd -d -f /etc/chrony/chrony.conf
autostart=true
autorestart=true
priority=10
startsecs=3
startretries=5
stdout_logfile=/var/log/timeserver/chrony.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/var/log/timeserver/chrony.log
stderr_logfile_maxbytes=5MB

; -----------------------------------------------------------------------
; 2. gpsd — czyta NMEA z urządzenia GPS USB
;    %(ENV_GPS_DEVICE)s pochodzi z entrypoint.sh (export GPS_DEVICE)
; -----------------------------------------------------------------------
[program:gpsd]
command=/usr/sbin/gpsd -N -n %(ENV_GPS_DEVICE)s
autostart=true
autorestart=true
priority=20
startsecs=2
startretries=10
stdout_logfile=/var/log/timeserver/gpsd.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/var/log/timeserver/gpsd.log
stderr_logfile_maxbytes=5MB

; -----------------------------------------------------------------------
; 3. gps-time-bridge — bridge: gpspipe (gpsd) → SHM 0 → chrony
;    Ma wbudowany retry loop co 10s gdy SHM jeszcze nie gotowy
;    startretries=10 daje czas na start chrony i gpsd
; -----------------------------------------------------------------------
[program:gps-time-bridge]
command=python3 /app/gps-time-bridge.py
autostart=true
autorestart=true
priority=30
startsecs=2
startretries=10
stdout_logfile=/var/log/timeserver/gps-time-bridge.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/var/log/timeserver/gps-time-bridge.log
stderr_logfile_maxbytes=5MB

; -----------------------------------------------------------------------
; 4. server.py — panel webowy HTTP
;    Zmienne środowiskowe przekazywane z docker run / docker-compose
; -----------------------------------------------------------------------
[program:timeserver]
command=python3 /app/server.py
autostart=true
autorestart=true
priority=40
startsecs=3
startretries=5
environment=HOST="0.0.0.0",PORT="%(ENV_PORT)s",ADMIN_USER="%(ENV_ADMIN_USER)s",ADMIN_PASS="%(ENV_ADMIN_PASS)s",SECRET_KEY="%(ENV_SECRET_KEY)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",LOG_DIR="/var/log/timeserver"
stdout_logfile=/var/log/timeserver/server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/timeserver/server.log
stderr_logfile_maxbytes=10MB
